# deploy/systemd/doktok-restore-import-ondemand.service
[Unit]
Description=DokTok on-demand portable restore (consumes the request file, then imports as root)
Documentation=https://github.com/lucianhanga/doktok-ng/blob/main/deploy/systemd/README.md
After=docker.service
Requires=docker.service
# Root-side cooldown: at most 1 restore per hour regardless of how many request files arrive. The
# backend is the first gate (confirm-to-destroy + validated-staged-id + single-flight); this is belt
# and braces in case a request file is dropped directly onto the box.
StartLimitIntervalSec=3600
StartLimitBurst=1

[Service]
Type=oneshot
WorkingDirectory=/opt/doktok
EnvironmentFile=/etc/doktok/backup.env
# (a) Read the staged_id out of the request file AND validate the uuid4-hex format (F-29, #641) -
# a malformed id fails this ExecStartPre, so the destructive ExecStart never runs - (b) DELETE the
# request file FIRST so a slow/failed restore can't re-trigger an endless loop, then (c) run the
# importer under a flock so two restores (or a restore + a backup timer) can never overlap.
# DOKTOK_BACKUP_DIR comes from the env file; the default matches the .path unit above. The
# destructive importer takes its own safety snapshot first.
ExecStartPre=/bin/sh -c 'sed -n "s/.*\\"staged_id\\":[[:space:]]*\\"\\([^\\"]*\\)\\".*/\\1/p" "${DOKTOK_BACKUP_DIR}/status/requests/restore.request" | grep -E "^[0-9a-f]{32}$" > /run/doktok-restore-staged-id'
ExecStartPre=/bin/rm -f ${DOKTOK_BACKUP_DIR}/status/requests/restore.request
ExecStart=/bin/sh -c '/usr/bin/flock -n /run/doktok-restore-import.lock /opt/doktok/deploy/restore-import.sh "$(cat /run/doktok-restore-staged-id)"'
# A restore is heavy and IO-bound; let it use the box but stay below interactive starvation.
Nice=10
IOSchedulingClass=best-effort
CPUQuota=200%
MemoryMax=4G
